# devcontainer/build/install-ops-tools/29999-openvpn/supervisord.conf
[program:openvpn]
command = openvpn --config /root/.config/openvpn/vpn.ovpn --auth-user-pass /root/.config/openvpn/auth.txt
user = root
autostart = false
autorestart = true
stdout_logfile=/var/log/supervisor/%(program_name)s/stdout.log
stderr_logfile=/var/log/supervisor/%(program_name)s/stderr.log
pidfile=/var/run/%(program_name)s.pid
