# devcontainer/build/install-base-tools/00007-cron/supervisord.conf
[program:cron]
command = /usr/sbin/cron -f
user = root
autostart = false
autorestart = true
stdout_logfile=/var/log/supervisor/%(program_name)s/stdout.log
stderr_logfile=/var/log/supervisor/%(program_name)s/stderr.log
pidfile=/var/run/%(program_name)s.pid