# devcontainer/build/install-base-tools/00006-ssh/supervisord.conf
[program:sshd]
command = /usr/sbin/sshd -D
user = root
autostart = true
autorestart = true
priority = 1
stdout_logfile=/var/log/supervisor/%(program_name)s/stdout.log
stderr_logfile=/var/log/supervisor/%(program_name)s/stderr.log
pidfile=/var/run/%(program_name)s.pid